# airootfs/etc/skel/.config/autostart/spicetheme.desktop
[Desktop Entry]
Name=Spice Theme
Exec=spicetheme
Icon=utilities-terminal
Type=Application
Categories=System;Utility;